# supervisord.conf
[supervisord]
nodaemon=true

[program:shell]
command=/bin/bash
autostart=true
autorestart=true
